# powerpill.service
[Unit]
Description=Download packages with Powerpill

[Service]
ExecStart=powerpill -Syuw
